# scripts/automation/cron/healthcheck.cron
# Run system health check every 5 minutes
*/5 * * * * /home/healthcheck/scripts/health-check.sh
